# smart_test/crontab
0 0 * * * root /home/pauls/tidning_monitoring/smart_test/smart_test.sh
